# files/lib/systemd/user/niri.service.d/niri-session.conf
[Unit]
Wants=dms.service
Wants=mate-polkit.service
