# files/usr/lib/systemd/user/noctalia.service
[Unit]
Description=Noctalia Shell
After=niri.service

[Service]
Type=simple
ExecStart=/usr/bin/noctalia
Restart=on-failure
