# files/system/usr/lib/systemd/user/dms.service
[Unit]
Description=Dank Material Shell
After=niri.service

[Service]
Type=simple
ExecStart=/usr/bin/dms run
Restart=on-failure
