# files/lib/systemd/user/noctalia.service
[Unit]
Description=Noctalia Shell
After=niri.service

[Service]
Type=simple
ExecStart=/usr/bin/qs -c noctalia-shell
Restart=on-failure
